# private_dot_config/autostart/simpread.desktop
[Desktop Entry]
Name=simpread
Exec=/usr/bin/Simpread
Terminal=False
Type=Application
