# backup/autostart/private_简悦 · 同步助手-1.2.0.AppImage.desktop
[Desktop Entry]
Exec='/mine2/obsidian/work/mine/简悦 · 同步助手-1.2.0.AppImage'
Icon=
Name=简悦 · 同步助手-1.2.0.AppImage
Path=
Terminal=False
Type=Application
