# backup/autostart/private_simpleread.AppImage.desktop
[Desktop Entry]
Exec=/mine/obsidian/obsidian@localrestapi/simpleread.AppImage
Icon=
Name=simpleread.AppImage
Path=
Terminal=False
Type=Application
